# com.resonanceenergy.ncl-brain.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.resonanceenergy.ncl-brain</string>

    <!-- Wrapper script sources ~/.env and the project .env before exec-ing uvicorn,
         ensuring all API keys and runtime variables are present. launchd does not
         run a login shell so .bashrc/.zshenv are never sourced without this wrapper. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/dev/NCL/scripts/launch-brain.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__/dev/NCL</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/dev/NCL/logs/ncl-brain-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/dev/NCL/logs/ncl-brain-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>PYTHONPATH</key>
        <string>__HOME__/dev/NCL</string>
    </dict>
</dict>
</plist>
