# config/com.resonanceenergy.ncl-orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Strike Point Orchestrator — mandate watcher service -->
    <!-- Install: cp config/com.resonanceenergy.ncl-orchestrator.plist ~/Library/LaunchAgents/ -->
    <!-- Load:    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.resonanceenergy.ncl-orchestrator.plist -->
    <!-- Unload:  launchctl bootout gui/$(id -u)/com.resonanceenergy.ncl-orchestrator -->

    <key>Label</key>
    <string>com.resonanceenergy.ncl-orchestrator</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.pyenv/shims/python3</string>
        <string>-m</string>
        <string>runtime.strike_point_orchestrator</string>
        <string>--watch</string>
        <string>--poll</string>
        <string>30</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__/dev/NCL</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:__HOME__/.pyenv/shims</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PYTHONPATH</key>
        <string>__HOME__/dev/NCL</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/dev/NCL/logs/orchestrator-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/dev/NCL/logs/orchestrator-stderr.log</string>

    <key>RunAtLoad</key>
    <true/><!-- Start automatically on login -->

    <key>KeepAlive</key>
    <true/><!-- Restart if it crashes -->

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
