# config/com.resonanceenergy.ncl-logrotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.resonanceenergy.ncl-logrotate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/dev/NCL/scripts/rotate-logs.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__/dev/NCL</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/dev/NCL/logs/logrotate.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/dev/NCL/logs/logrotate-error.log</string>
</dict>
</plist>
